# chinlab/databases/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)
;chmod=0700                 ; socket file mode (default 0700)
;chown=nobody:nogroup       ; socket file uid:gid owner
;username=user              ; (default is no username (open server))
;password=123               ; (default is no password (open server))

[inet_http_server]
port=0.0.0.0:19999                  ; (ip_address:port specifier, *:port for all iface)
username=admin               ; (default is no username (open server))
password=mhealth365          ; (default is no password (open server))

[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)
;umask=022                   ; (process file creation umask;default 022)
;user=chrism                 ; (default is current user, required if root)
;identifier=supervisor       ; (supervisord identifier, default is 'supervisor')
;directory=/tmp              ; (default is not to cd during start)
;nocleanup=true              ; (don't clean up tempfiles at start;default false)
;childlogdir=/tmp            ; ('AUTO' child log dir, default $TEMP)
;environment=KEY="value"     ; (key value pairs to add to environment)
;strip_ansi=false            ; (strip ansi escape codes in logs; def. false)

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket
;serverurl=http://127.0.0.1:19999 ; use an http:// url to specify an inet socket
;username=chris              ; should be same as http_username if set
;password=123                ; should be same as http_password if set
;prompt=mysupervisor         ; cmd line prompt (default "supervisor")
;history_file=~/.sc_history  ; use readline history if available

[program:active_sms_phone_code]
command=php /home/wwwroot/rongyao/yii  rabbit/index active_sms_phone_code 
autorstart=true
process_name=%(program_name)s_%(process_num)02d
numprocs=2
stdout_logfile=/var/log/www/command/active_sms_phone_code.log


[program:sms_phone_code]
command=php /home/wwwroot/huobanys/yii  rabbit/index sms_phone_code
autorstart=true
process_name=%(program_name)s_%(process_num)02d
numprocs=2
stdout_logfile=/var/log/www/command/elasticsearch.log

[program:user_info_upinsert]
command=php /home/wwwroot/huobanys/yii  rabbit/index user_info_upinsert
autorstart=true
process_name=%(program_name)s_%(process_num)02d
numprocs=6
stdout_logfile=/var/log/www/command/user_info_upinsert.log

[program:system_cancel_order]
command=php /home/wwwroot/huobanys/yii rabbit/index system_cancel_order
autorstart=true
process_name=%(program_name)s_%(process_num)02d
numprocs=4
stdout_logfile=/var/log/www/command/system_cancel_order.log

[program:system_cancel_first_order]
command=php /home/wwwroot/huobanys/yii  rabbit/index system_cancel_first_order
autorstart=true
process_name=%(program_name)s_%(process_num)02d
numprocs=4
stdout_logfile=/var/log/www/command/system_cancel_first_order.log

[program:user_push_from_sys]
command=php /home/wwwroot/huobanys/yii rabbit/index user_push_from_sys 
autorstart=true
process_name=%(program_name)s_%(process_num)02d
numprocs=4
stdout_logfile=/var/log/www/command/user_push_from_sys.log

[program:user_sms_from_sys]
command=php /home/wwwroot/huobanys/yii rabbit/index user_sms_from_sys
autorstart=true
process_name=%(program_name)s_%(process_num)02d
numprocs=4
stdout_logfile=/var/log/www/command/user_sms_from_sys.log

[program:order_sms_expire_sys]
command=php /home/wwwroot/huobanys/yii  rabbit/index order_sms_expire_sys
autorstart=true
process_name=%(program_name)s_%(process_num)02d
numprocs=4
stdout_logfile=/var/log/www/command/order_sms_expire_sys.log

[program:order_sms_from_pay]
command=php /home/wwwroot/huobanys/yii  rabbit/index order_sms_from_pay
autorstart=true
process_name=%(program_name)s_%(process_num)02d
numprocs=4
stdout_logfile=/var/log/www/command/order_sms_from_pay.log

[program:news_one_cache]
command=php /home/wwwroot/huobanys/yii  rabbit/index news_one_cache
autorstart=true
process_name=%(program_name)s_%(process_num)02d
numprocs=4
stdout_logfile=/var/log/www/command/news_one_cache.log

[program:news_more_cache]
command=php /home/wwwroot/huobanys/yii  rabbit/index news_more_cache
autorstart=true
process_name=%(program_name)s_%(process_num)02d
numprocs=4
stdout_logfile=/var/log/www/command/news_more_cache.log

[program:news_push_info]
command=php /home/wwwroot/huobanys/yii  rabbit/index news_push_info
autorstart=true
process_name=%(program_name)s_%(process_num)02d
numprocs=4
stdout_logfile=/var/log/www/command/news_push_info.log

[program:news_show_end]
command=php /home/wwwroot/huobanys/yii  rabbit/index news_show_end
autorstart=true
process_name=%(program_name)s_%(process_num)02d
numprocs=4
stdout_logfile=/var/log/www/command/news_show_end.log

[program:news_card_create]
command=php /home/wwwroot/huobanys/yii  rabbit/index news_card_create
autorstart=true
process_name=%(program_name)s_%(process_num)02d
numprocs=4
stdout_logfile=/var/log/www/command/news_card_create.log

[program:goods_info_cache]
command=php /home/wwwroot/huobanys/yii  rabbit/index goods_info_cache
autorstart=true
process_name=%(program_name)s_%(process_num)02d
numprocs=4
stdout_logfile=/var/log/www/command/goods_info_cache.log


[program:system_refund_order]
command=php /home/wwwroot/huobanys/yii  rabbit/index system_refund_order
autorstart=true
process_name=%(program_name)s_%(process_num)02d
numprocs=2
stdout_logfile=/var/log/www/command/system_refund_order.log


;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;

[program:b_sms_phone_code]
command=php /home/wwwroot/before.huobanys/yii  rabbit/index sms_phone_code
autorstart=true
process_name=%(program_name)s_%(process_num)02d
numprocs=1
stdout_logfile=/var/log/www/command/b_elasticsearch.log

[program:b_user_info_upinsert]
command=php /home/wwwroot/before.huobanys/yii  rabbit/index user_info_upinsert
autorstart=true
process_name=%(program_name)s_%(process_num)02d
numprocs=1
stdout_logfile=/var/log/www/command/b_user_info_upinsert.log

[program:b_system_cancel_order]
command=php /home/wwwroot/before.huobanys/yii rabbit/index system_cancel_order
autorstart=true
process_name=%(program_name)s_%(process_num)02d
numprocs=1
stdout_logfile=/var/log/www/command/b_system_cancel_order.log

[program:b_system_cancel_first_order]
command=php /home/wwwroot/before.huobanys/yii  rabbit/index system_cancel_first_order
autorstart=true
process_name=%(program_name)s_%(process_num)02d
numprocs=1
stdout_logfile=/var/log/www/command/b_system_cancel_first_order.log

[program:b_user_push_from_sys]
command=php /home/wwwroot/before.huobanys/yii rabbit/index user_push_from_sys 
autorstart=true
process_name=%(program_name)s_%(process_num)02d
numprocs=1
stdout_logfile=/var/log/www/command/b_user_push_from_sys.log


[program:b_user_sms_from_sys]
command=php /home/wwwroot/before.huobanys/yii rabbit/index user_sms_from_sys
autorstart=true
process_name=%(program_name)s_%(process_num)02d
numprocs=1
stdout_logfile=/var/log/www/command/b_user_sms_from_sys.log

[program:b_order_sms_expire_sys]
command=php /home/wwwroot/before.huobanys/yii  rabbit/index order_sms_expire_sys
autorstart=true
process_name=%(program_name)s_%(process_num)02d
numprocs=1
stdout_logfile=/var/log/www/command/b_order_sms_expire_sys.log

[program:b_order_sms_from_pay]
command=php /home/wwwroot/before.huobanys/yii  rabbit/index order_sms_from_pay
autorstart=true
process_name=%(program_name)s_%(process_num)02d
numprocs=1
stdout_logfile=/var/log/www/command/b_order_sms_from_pay.log

[program:b_news_one_cache]
command=php /home/wwwroot/before.huobanys/yii  rabbit/index news_one_cache
autorstart=true
process_name=%(program_name)s_%(process_num)02d
numprocs=1
stdout_logfile=/var/log/www/command/b_news_one_cache.log

[program:b_news_more_cache]
command=php /home/wwwroot/before.huobanys/yii  rabbit/index news_more_cache
autorstart=true
process_name=%(program_name)s_%(process_num)02d
numprocs=1
stdout_logfile=/var/log/www/command/b_news_more_cache.log

[program:b_news_push_info]
command=php /home/wwwroot/before.huobanys/yii  rabbit/index news_push_info
autorstart=true
process_name=%(program_name)s_%(process_num)02d
numprocs=1
stdout_logfile=/var/log/www/command/b_news_push_info.log

[program:b_news_show_end]
command=php /home/wwwroot/before.huobanys/yii  rabbit/index news_show_end
autorstart=true
process_name=%(program_name)s_%(process_num)02d
numprocs=1
stdout_logfile=/var/log/www/command/b_news_show_end.log


[program:b_news_card_create]
command=php /home/wwwroot/before.huobanys/yii  rabbit/index news_card_create
autorstart=true
process_name=%(program_name)s_%(process_num)02d
numprocs=1
stdout_logfile=/var/log/www/command/b_news_card_create.log

[program:b_goods_info_cache]
command=php /home/wwwroot/before.huobanys/yii  rabbit/index goods_info_cache
autorstart=true
process_name=%(program_name)s_%(process_num)02d
numprocs=1
stdout_logfile=/var/log/www/command/b_goods_info_cache.log


[program:b_system_refund_order]
command=php /home/wwwroot/before.huobanys/yii  rabbit/index system_refund_order
autorstart=true
process_name=%(program_name)s_%(process_num)02d
numprocs=1
stdout_logfile=/var/log/www/command/b_system_refund_order.log


